# scripts/launchd/com.hotcollect.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.hotcollect.server</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd /Users/mima0000/Desktop/hot_collect &amp;&amp; mkdir -p .data/logs &amp;&amp; exec uvicorn server:app --host 127.0.0.1 --port 8000</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/mima0000/Desktop/hot_collect</string>

    <key>StandardOutPath</key>
    <string>/Users/mima0000/Desktop/hot_collect/.data/logs/server.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mima0000/Desktop/hot_collect/.data/logs/server.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>
  </dict>
</plist>

